# deploy/curfewbot.service
# CurfewBot systemd service — non-Docker alternative
#
# Installation:
#   sudo cp deploy/curfewbot.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable curfewbot
#   sudo systemctl start curfewbot
#
# Logs:
#   sudo journalctl -u curfewbot -f

[Unit]
Description=CurfewBot Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/CurfewBot
EnvironmentFile=/home/ec2-user/CurfewBot/.env
Environment=DB_DIR=/home/ec2-user/CurfewBot/data
ExecStartPre=/bin/mkdir -p /home/ec2-user/CurfewBot/data
ExecStart=/usr/bin/python3 src/curfewbot.py
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/ec2-user/CurfewBot/data

[Install]
WantedBy=multi-user.target
